# yoctopuce/yVirtualHub
#!/sbin/openrc-run

command="/usr/sbin/VirtualHub"
